# roles/download_backups/templates/rclone.cron
0 6 * * * /usr/bin/rclone --config /etc/rclone.conf sync -l /var/www/html/ s3:gluster.org-download-backups 
